# productai.service
[Unit]
Description=ProductAI - AI-powered Product Management
After=network.target

[Service]
Type=simple
User=adrian
WorkingDirectory=/opt/services/productai
Environment=BASE_PATH=/productai
ExecStart=/usr/local/bin/uv run uvicorn productai.app:app --host 127.0.0.1 --port 5010
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
